# systemd/quest-dashboard.service
[Unit]
Description=Quest Dashboard — static HTML + tag-mutation API server
Documentation=https://github.com/ytrofr/claude-code-quest
After=network.target

[Service]
Type=simple
Environment=QUEST_HOST=127.0.0.1
Environment=QUEST_PORT=8770
# server.py: serves static site/, plus POST /api/tags/{add,remove,bind,unbind}
# for UI tag mutations. Renders on startup + after every mutation.
ExecStart=/usr/bin/python3 %h/.claude/skills/quest/server.py
Restart=on-failure
RestartSec=5s

# Hardening — port-bound to localhost only via QUEST_HOST above
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/.claude/quest

[Install]
WantedBy=default.target
